# automation/com.jobsearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd plist for daily job search runs.
  UPDATE all paths marked with /path/to/your/project before installing.

  Install:
    cp automation/com.jobsearch.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.jobsearch.plist

  For Linux, use cron instead:
    0 7 * * * /path/to/your/project/automation/run_search.sh >> /path/to/your/project/automation/logs/cron.log 2>&1

  For Windows (Task Scheduler), create a task that runs:
    bash /path/to/your/project/automation/run_search.sh
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobsearch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- UPDATE: absolute path to run_search.sh -->
        <string>/path/to/your/project/automation/run_search.sh</string>
    </array>

    <!-- UPDATE: absolute path to project root -->
    <key>WorkingDirectory</key>
    <string>/path/to/your/project</string>

    <!-- Run daily at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- UPDATE: absolute paths for log output -->
    <key>StandardOutPath</key>
    <string>/path/to/your/project/automation/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/your/project/automation/logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
