# automation/com.jobsearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd plist for daily job search runs.
  UPDATE all paths marked with /path/to/your/project before installing.

  Install:
    cp automation/com.jobsearch.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.jobsearch.plist

  IMPORTANT (macOS laptops): launchd fires scheduled jobs during DarkWake, where
  the Mac will not stay awake long enough for the long evaluation API call and it
  times out. Two mitigations are baked in below — `caffeinate` keeps the Mac awake
  for the run, and RunAtLoad catches up a run missed while asleep/off. For full
  reliability, also schedule a power-on wake a couple minutes before the run (needs
  sudo, one-time):
    sudo pmset repeat wakeorpoweron MTWRFSU 06:58:00   # match to your run time

  For Linux, use cron instead:
    0 7 * * * /path/to/your/project/automation/run_search.sh >> /path/to/your/project/automation/logs/cron.log 2>&1

  For Windows (Task Scheduler), create a task that runs:
    bash /path/to/your/project/automation/run_search.sh
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobsearch</string>

    <!-- caffeinate keeps the Mac fully awake for the whole run. launchd fires this
         job during DarkWake/maintenance-sleep, which otherwise stalls the long
         evaluation API call. -i no idle sleep, -m no disk sleep, -s no system sleep (AC). -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-m</string>
        <string>-s</string>
        <string>/bin/bash</string>
        <!-- UPDATE: absolute path to run_search.sh -->
        <string>/path/to/your/project/automation/run_search.sh</string>
    </array>

    <!-- UPDATE: absolute path to project root -->
    <key>WorkingDirectory</key>
    <string>/path/to/your/project</string>

    <!-- Run daily at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- PATH should include wherever `claude` is installed.
         Claude Code's native installer puts it at ~/.local/bin/claude — plist files do not
         expand ~, so UPDATE the absolute path below to your home directory if applicable.
         Run `which claude` to confirm. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- UPDATE: absolute paths for log output -->
    <key>StandardOutPath</key>
    <string>/path/to/your/project/automation/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/your/project/automation/logs/launchd_stderr.log</string>

    <!-- Run on load (login/boot) to catch up if the scheduled run was missed while
         the Mac was asleep/off. run_search.sh's catch-up guard exits early if today's
         report already exists, so this won't duplicate a successful scheduled run. -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
