# automation/com.jobsearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobsearch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOURNAME/code/jobs/automation/run_search.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURNAME/code/jobs</string>

    <!-- Run daily at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/YOURNAME</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/YOURNAME/code/jobs/automation/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURNAME/code/jobs/automation/logs/launchd_stderr.log</string>

    <!-- Run on load if we missed the scheduled time (machine was asleep) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
